# dev/provision/crontab
*/5 * * * * php /vagrant/bin/magento cron:run
